# systems/lenovo/x220/overlay/etc/udev/rules.d/95-keymap.rules
SUBSYSTEM=="input", ATTRS{name}=="AT Translated Set 2 keyboard", RUN+="keymap $name /lib/udev/keymaps/lenovo-thinkpad_x220_tablet"
